# start/deploy/product/supervisor_start.ini
[supervisord]
daemon=true
[program:nginx]
command=nginx
autostart=true
autorestart=true
startretries=3
#关闭所有进程
stopasgroup=true
# 同上
killasgroup = true

[program:start_uwsgi]
directory=/opt/start
environment=START_ENV_TYPE="product"
command=uwsgi --ini product_uwsgi.ini
numprocs=1
autostart=true
autorestart=true
startretries=3
#关闭所有进程
stopasgroup=true
# 同上
killasgroup = true

[program:celery_work]
directory=/opt/start
environment=START_ENV_TYPE="product"
command=Celery -A demo.demo_task worker -l INFO
#numprocs=1
autostart=true
autorestart=true
startretries=3
#关闭所有进程
stopasgroup=true
# 同上
killasgroup = true
redirect_stderr = true
stdout_logfile_backups = 5
stderr_logfile = /opt/start/logs/worker_stderr.log
stdout_logfile = /opt/start/logs/worker_stdout.log

[program:celery_beat]
directory=/opt/start
environment=START_ENV_TYPE="product"
command=celery -A demo.demo_job beat -l INFO --scheduler django_celery_beat.schedulers:DatabaseScheduler
autostart=true
autorestart=true
startretries=3
#关闭所有进程
stopasgroup=true
# 同上
killasgroup = true
redirect_stderr = true
stdout_logfile_backups = 5
stderr_logfile = /opt/start/logs/beat_stderr.log
stdout_logfile = /opt/start/logs/beat_stdout.log